# cron.hourly/pktlog-clean
#!/bin/sh

find /data/log/. -type f -name '*.dat' -mtime +3 | xargs rm -f
